# deploy/supervisor-raautomotriz-worker.conf
; Worker de cola para RA Automotriz (Laravel)
; Copiar a: /etc/supervisor/conf.d/raautomotriz-worker.conf
; Luego: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start raautomotriz-worker

[program:raautomotriz-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/raautomotriz/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
directory=/var/www/raautomotriz
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/raautomotriz/storage/logs/worker.log
stopwaitsecs=60
